# profile.d/alias.sh
#!/bin/ksh

alias ll='ls -lart'
alias zshconfig='gvim ~/.zshrc'
alias envconfig='gvim ~/env.sh'
alias plantuml='java -jar ~/.m2/repository/net/sourceforge/plantuml/plantuml/8018/plantuml-8018.jar'

