# profile.d/alias.sh
alias ll='ls -lart'
alias plantuml='java -jar ~/.m2/repository/net/sourceforge/plantuml/plantuml/8018/plantuml-8018.jar'
