# profile.d/alias.sh
#!/bin/ksh

alias ll='ls -lart'
alias zshconfig='gvim ~/.zshrc'
alias envconfig='gvim ~/env.sh'
alias sbt_eclipse="sbt 'eclipse with-source=true'"
alias plantuml='java -jar ~/.m2/repository/net/sourceforge/plantuml/plantuml/8018/plantuml-8018.jar'
alias ws='cd ~/development/workspace' # ws => workspace
alias fix_vpn='docker-machine scp /etc/resolv.conf default:/etc/resolv.conf'
