# profile.d/variables.sh
export LC_ALL='en_US.UTF-8'
export DISPLAY='10.0.2.2:0' # host ip default (virtualbox)
